# deploy/taskwarrior-web.service.tmpl
[Unit]
Description=taskwarrior-web-portal - local web UI for Taskwarrior
Documentation=https://github.com/furan917/taskwarrior-web-portal
# Light dependency on basic.target (the user manager's "ready for services"
# goal). default.target IS the synthetic that pulls things in via WantedBy,
# so After=default.target reads as "after our own goal" - backwards.
After=basic.target

[Service]
Type=simple
ExecStart=__BIN__
# Restart on crash, not on clean exit. Mirrors the launchd KeepAlive dict
# (SuccessfulExit=false, Crashed=true) so a clean `systemctl --user stop`
# stays stopped.
Restart=on-failure
RestartSec=10
TimeoutStopSec=10

# stdout/stderr go to the systemd journal. Tail with:
#   journalctl --user -u taskwarrior-web-portal -f
# These are the panic safety net only - the binary writes its primary
# structured slog stream to ${XDG_STATE_HOME:-~/.local/state}/
# taskwarrior-web-portal/app.log with size-based rotation (10 MB/file, 3 backups,
# 30 day max age, gzip-compressed) via lumberjack. The journal stream
# should remain near-empty under normal operation.
StandardOutput=journal
StandardError=journal
StandardInput=null

# Discoverability for `task`. %h is systemd's home-dir specifier. Covers the
# common install locations: ~/.local/bin (pip --user, uv tool, custom build),
# Linuxbrew default prefix, /usr/local (manual builds), distro packages
# (/usr/bin). If your `task` lives elsewhere, override here or via a drop-in
# at ~/.config/systemd/user/taskwarrior-web-portal.service.d/path.conf.
Environment=PATH=%h/.local/bin:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin

# XDG_STATE_HOME is pinned by install.sh at install time so the install
# script and the running binary can't disagree on the log directory. If
# the user's interactive shell sets a non-default XDG_STATE_HOME (e.g. in
# ~/.zshrc), install.sh sees it and bakes the resolved value here;
# systemd's --user manager doesn't auto-import shell rc-file env, so
# without this the binary would fall back to ~/.local/state and the
# install-time `chmod 700` would have happened on a different directory.
Environment=XDG_STATE_HOME=__XDG_STATE_HOME__

# Light hardening. The binary is single-user, local-loopback only - heavy
# sandboxing (ProtectSystem=strict, PrivateUsers, etc.) tends to break
# Taskwarrior's hooks and homedir access without buying meaningful
# protection in this threat model. NoNewPrivileges is the cheap win.
NoNewPrivileges=true

[Install]
# default.target is the user-session conventional goal. Linger off (the
# default) means the service starts at first login of this user; turn on
# `loginctl enable-linger $USER` if you want it running before any login.
WantedBy=default.target
